# cnc-pinball.service
[Unit]
Description=CnC Pinball GUI
# Varjunk meg, amig a multi-user celallapot elerheto (boot vege),
# de NE varjunk grafikus retegre (nincs is X/Wayland telepitve - a
# GUI kozvetlenul DRM/KMS-re rajzol, lasd score_gui.py kmsdrm-detektalas).
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=gonzodr
Group=gonzodr

# A felhasznalo legyen tagja a 'video', 'render' es 'dialout' csoportoknak,
# mert kell: DRM/KMS irasi jog (video, render) es soros port eleres (dialout).
# Ellenorzes: groups gonzodr
# Hozzaadas, ha hianyzik: sudo usermod -aG video,render,dialout gonzodr

WorkingDirectory=/home/gonzodr/CnC-Pinball-GUI-V3/src
ExecStart=/home/gonzodr/CnC-Pinball-GUI-V3/venv/bin/python3 /home/gonzodr/CnC-Pinball-GUI-V3/src/main.py

# Ha a program elszall (pl. soros hiba, mpv crash), systemd inditsa ujra.
Restart=on-failure
RestartSec=3

# Naplozas a journal-ba, hogy `journalctl -u cnc-pinball -f`-fel kovetheto legyen
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
